# services/crosscompute-teamwork.service
[Service]
Type=oneshot
ExecStart=/root/Projects/crosscompute-teamwork/scripts/build.sh
